# packages/base/all/initrds/loader-initrd-files/src/etc/inittab
# Run /bin/sysinit without a controlling tty so that job control
# cannot be used to circumvent password security
::sysinit:/bin/sysinit

# Try booting automatically, start a shell if user interrupts it
::wait:-/bin/autoboot
::wait:-/bin/login
::wait:/bin/umount -a -r
::wait:/sbin/reboot -f

# Try to switch to rootfs if init receives SIGHUP or SIGQUIT
::restart:/bin/switchroot
