# packages/base/all/vendor-config-onl/src/boot.d/51.onl-platform-baseconf
#!/usr/bin/python

### BEGIN INIT INFO
# Provides:        onl-platform-baseconf
# Required-Start:
# Required-Stop:
# Default-Start:   S
# Default-Stop:
# Short-Description: Set up ONL Platform
### END INIT INFO

from onl.platform.baseconfig import baseconfig

baseconfig()



